# recipes-core/efi/files/boot.mount
[Unit]
Before=local-fs.target

[Mount]
What=/dev/disk/by-partlabel/efi
Where=/boot
Type=vfat
Options=noatime,sync
