# recipes-core/overlayfs/files/mnt-overlay.mount
[Unit]
Before=local-fs.target

[Mount]
What=/dev/disk/by-partlabel/overlay
Where=/mnt/overlay
Type=ext4
Options=inlinecrypt