# recipes-bsp/usb/files/usb_bind.rules
SUBSYSTEM=="udc", ACTION=="add", RUN+="/usr/bin/qusb rebind"
